# deploy/crontab
*/5 * * * * %(user)s %(manage)s poll_rss http://www.canberratimes.com.au/rssheadlines/federal-politics/article/rss.xml --follow
